# docker/supervisor_test/supervisord.conf
[program:nginx]
command=nginx

[program:nodeserver]
command=node /provision/static/server.js


[supervisord]
